# deploy/com.choiyounggi.kds-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.choiyounggi.kds-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/choeyeonggi/korea-data-suite/scripts/sync.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/choeyeonggi/korea-data-suite</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/choeyeonggi/Library/Logs/kds/sync.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/choeyeonggi/Library/Logs/kds/sync.err.log</string>
</dict>
</plist>
